# scripts/backup-config.timer
[Unit]
Description=Mediastack: run configuration backup weekly

[Timer]
# Weekly, mid-morning Sunday local time. Persistent=true runs a missed backup
# after boot if the server was off at the scheduled time.
OnCalendar=Sun *-*-* 10:15:00
RandomizedDelaySec=15min
Persistent=true
Unit=backup-config.service

[Install]
WantedBy=timers.target
